# deploy/supervisor/horizon.conf
[program:horizon]
process_name=%(program_name)s
command=php /var/www/hackatonshik/artisan horizon
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/hackatonshik/storage/logs/horizon.log
stopwaitsecs=3600
